# rpi/youtube-stream-service/etc/systemd/system/youtube.socket
[Unit]
Description=The fifo pipe for streaming youtube

[Socket]
ListenFIFO=/run/youtube

[Install]
WantedBy=sockets.target
